# templates/systemd/board-manager-v4.service
[Unit]
Description=Board LAN Hub (IPv4)
After=network.target

[Service]
Type=simple
WorkingDirectory={{APPDIR}}
ExecStart={{PYTHON}} -m uvicorn app.main:app --host 0.0.0.0 --port {{APIPORT}}
Restart=always

[Install]
WantedBy=multi-user.target
